# config/headnode/rc.local
#!/bin/sh

# Rules for routing traffic from the nodes out
iptables-restore < /etc/iptables.rules

# Limit the GPU clocks
nvidia-smi -pm 1
nvidia-smi -lgc 1300,1300

# Set CPUs to performance mode
for i in `seq 0 1 63`; do
        echo performance > /sys/devices/system/cpu/cpu${i}/cpufreq/scaling_governor
done
