# config/headnode/rc.local
#!/bin/bash

# Rules for routing traffic from the nodes out
iptables-restore < /etc/iptables.rules

# Limit the GPU clocks
nvidia-smi -lgc 1200,1200

exit 0
